# PetaLinux-SDK/project-spec/meta-user/recipes-apps/bootscript/files/bootscript.service
[Unit]
Description=simple-script
  
[Service]
ExecStart=/usr/bin/bootscript.sh
StandardOutput=journal+console
  
[Install]
WantedBy=multi-user.target